# config/includes.chroot/etc/xdg/autostart/trust-desktop-files.desktop
[Desktop Entry]
Type=Application
Name=Trust Desktop Files
Comment=Mark desktop shortcuts as trusted
Exec=/usr/local/bin/trust-desktop-files
NoDisplay=true
X-GNOME-Autostart-enabled=true
